# scripts/hummingbird.service
[Unit]
Description=Backyard Hummers - Hummingbird Feeder Camera
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=admin
WorkingDirectory=/home/admin/LocalHummingBirdCam
EnvironmentFile=/home/admin/LocalHummingBirdCam/.env
ExecStart=/home/admin/LocalHummingBirdCam/venv/bin/python main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
