# scripts/hummingbird.service
[Unit]
Description=Backyard Hummers - Hummingbird Feeder Camera
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/LocalHummingBirdCam
EnvironmentFile=/home/pi/LocalHummingBirdCam/.env
ExecStart=/home/pi/LocalHummingBirdCam/venv/bin/python main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
