# scripts/hummingbird-updater.service
[Unit]
Description=Backyard Hummers - Auto Update Check
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=admin
ExecStart=/home/admin/LocalHummingBirdCam/scripts/auto_update.sh
